# crontab
* * * * * echo "Current container's time:"$(date +"%Y-%m-%d %H:%M")
